# ops/systemd/homework-cert-renew.timer
[Unit]
Description=Run homework TLS certificate renewal twice daily

[Timer]
OnCalendar=*-*-* 03,15:00:00
Persistent=true
Unit=homework-cert-renew.service

[Install]
WantedBy=timers.target
